# test-components/postgresql/docker-centos/postgres.ini
[program:postgres]
command=bash -c "/usr/pgsql/bin/postgres -D /var/lib/pgsql/data -p 5432"
process_name=%(process_num)02d
user=postgres
childlogdir=/var/log/postgres
stdout_logfile=/var/log/postgres/%(program_name)s-stdout.log
stderr_logfile=/var/log/postgres/%(program_name)s-stderr.log
stdout_logfile_maxbytes=2GB
stderr_logfile_maxbytes=2GB
stdout_logfile_backups=10
stderr_logfile_backups=10
redirect_stderr=false
autorestart=true
stopsignal=TERM
stopasgroup=true
exitcodes=1
